# powermeter.service
[Unit]
Description=my Powermeter

[Service]
ExecStart=/var/log/powermeter.sh start

[Install]
WantedBy=multi-user.target