# powermeter.service
[Unit]
Description=my Powermeter

[Service]
ExecStart=/home/pi/powermeter.sh start

[Install]
WantedBy=multi-user.target
